# doc/deploy/supervisor.ini
[program:nginx]
directory=/opt
command=nginx -c /etc/nginx/nginx.conf
autostart=true
stopwaitsecs=3600

[program:run]
directory=/opt/OA
command=/usr/local/python3/bin/gunicorn -k gevent -b 0.0.0.0:9000 -w 1 run:app
autostart=true
stopwaitsecs=3600
